# ansible-docker-registry/roles/configure-docker-registry/files/docker-registry.service
[Unit]
Description=Docker Registry

[Service]
WorkingDirectory=/docker-registry
ExecStart=/usr/bin/docker-compose up
Restart=always

[Install]
WantedBy=multi-user.target
